# deploy/systemd/xspacedownloader-bg.service
[Unit]
Description=XSpace Downloader Background Downloader
After=network.target mysql.service xspacedownloader.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/xspacedownloader
Environment="PATH=/var/www/xspacedownloader/venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/var/www/xspacedownloader"
ExecStart=/var/www/xspacedownloader/venv/bin/python /var/www/xspacedownloader/bg_downloader.py
Restart=always
RestartSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/xspacedownloader/downloads /var/www/xspacedownloader/logs

# Logging
StandardOutput=append:/var/www/xspacedownloader/logs/bg-downloader.log
StandardError=append:/var/www/xspacedownloader/logs/bg-downloader-error.log

[Install]
WantedBy=multi-user.target